# solar-wallpaper.service
[Unit]
Description=Solar-based Dynamic Wallpaper Changer

[Service]
Type=oneshot
ExecStart=%h/.local/bin/solar-wallpaper.sh
Environment=DISPLAY=:0
